# units/idle-shutdown.timer
[Unit]
Description=Run Idle Shutdown Check Regularly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
